# .config/cron/crontab
# vi: set syntax=crontab:

0 22 * * * sh ~/.dotfiles/scripts/rgb/dim.sh
